# restic-backup.service
[Unit]
Description=Backup Files with Restic

[Service]
ExecStart=/usr/local/bin/restic-backup.sh

[Install]
WantedBy=multi-user.target
